# update-schedules-on-pioneer.cron
5,20,35,50 * * * *   nornetpp   if [ ! -e /nne/experiments/nne-node-experiments ] ; then git clone -q git@github.com:simula/nne-node-experiments.git /nne/experiments/nne-node-experiments ; else cd /nne/experiments/nne-node-experiments ] ; fi && git pull -q && cd /nne/experiments && /home/nornetpp/src/nne-scheduler-configurator/update-schedules && git push -f
