# update-schedules-on-pioneer.cron
*/15 * * * *   root   if [ ! -e /nne/nne-node-experiments/ ] ; then git clone -q https://github.com/simula/nne-node-experiments.git /nne/nne-node-experiments ; else cd /nne/nne-node-experiments/ && git pull -q && cd /nne/nne-node-experiments/ && ~nornetpp/src/nne-schedule ; fi
